# supervisor.ini
[program:food.hypertexthero.com]
command=/home/hth/webapps/foodnews/fn-env/bin/python2.7 /home/hth/webapps/foodnews/fn-env/bin/gunicorn_django --daemon -b 127.0.0.1:17734 -w 2 --max-requests 500
directory=/home/hth/webapps/foodnews/fn-env/fn
user=hth
autostart=true
autorestart=true
redirect_stderr=True